# src-tauri/crates/kriya-aggregator/packaging/box/kriyad.service
[Unit]
Description=kriyad — kriya on-prem evidence aggregator (single-tenant, no egress)
Documentation=https://github.com/sandeepshekhar26/kriya
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/kriyad
EnvironmentFile=/etc/kriyad/kriyad.env
# The license + mTLS material live in root-owned /etc/kriyad, which the DynamicUser sandbox
# (correctly) cannot read — PROVEN on a real host: the service died with EACCES on the license
# (CI kriyad-release run 1, 2026-07-02). LoadCredential= is the systemd-native fix: PID 1 reads
# them as root and exposes read-only copies to the sandboxed service under
# $CREDENTIALS_DIRECTORY (%d). The three cert files are loaded INDIVIDUALLY (run 2 proved a
# directory credential flattens to `ca_<file>`, not `ca/<file>`), under exactly the filenames
# kriyad expects, so KRIYAD_CA_DIR is the credentials root itself. Consequence: under systemd
# the CA files are REQUIRED at start (bootstrap with kriyd-ca.sh before `systemctl start`) —
# a production box must never silently fall back to plain HTTP anyway; the plain-HTTP dev mode
# remains available in manual runs. Needs systemd ≥ 247 — Ubuntu 22.04+, Debian 12+, RHEL 9+.
# On older hosts fall back to a static `kriyad` system user owning the files.
LoadCredential=license:/etc/kriyad/kriyad-license.json
LoadCredential=server.pem:/etc/kriyad/ca/server.pem
LoadCredential=server.key:/etc/kriyad/ca/server.key
LoadCredential=ca.pem:/etc/kriyad/ca/ca.pem
Environment=KRIYAD_LICENSE=%d/license
Environment=KRIYAD_CA_DIR=%d
Restart=on-failure
RestartSec=2

# --- Hardening -------------------------------------------------------------
# kriyad handles signed evidence and makes NO outbound calls, so it can run with
# almost nothing: a dynamic (transient) user, a read-only filesystem, no caps,
# and a syscall allowlist. Binds :8443 (>1024), so no privileged-port capability
# is needed and CapabilityBoundingSet can be empty.
DynamicUser=yes
StateDirectory=kriyad
# ^ SQLite store lands in /var/lib/kriyad, owned by the transient user.
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
ProtectProc=invisible
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
SystemCallArchitectures=native
UMask=0077

[Install]
WantedBy=multi-user.target
